# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g 'daemon off;'

[program:run_stream]
command=/run_stream.sh
